# darwin-wayland/userland/rootfs-skel/Library/LaunchDaemons/org.lunaos.compositor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.lunaos.compositor.plist — LaunchDaemon for luna-compositor
  Location: /Library/LaunchDaemons/org.lunaos.compositor.plist

  Starts luna-compositor as the lunaos user after:
    - seatd-darwin is running (/run/seatd.sock exists)
    - darwin-evdev-bridge is running (/dev/input/ has nodes)
    - IODRMShim.kext is loaded (/dev/dri/card0 exists)

  If the compositor crashes, launchd restarts it after 2 seconds.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.lunaos.compositor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/luna-compositor</string>
    </array>

    <key>UserName</key>
    <string>lunaos</string>

    <key>GroupName</key>
    <string>video</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Only keep alive if the required sockets/devices exist -->
        <key>PathState</key>
        <dict>
            <key>/run/seatd.sock</key>
            <true/>
            <key>/dev/dri/card0</key>
            <true/>
        </dict>
    </dict>

    <key>ThrottleInterval</key>
    <integer>2</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WAYLAND_DISPLAY</key>
        <string>wayland-0</string>
        <key>XDG_RUNTIME_DIR</key>
        <string>/run/user/501</string>
        <key>WLR_RENDERER</key>
        <string>pixman</string>
        <key>WLR_BACKENDS</key>
        <string>drm</string>
        <key>WLR_NO_HARDWARE_CURSORS</key>
        <string>1</string>
        <key>LIBSEAT_BACKEND</key>
        <string>seatd</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/luna-compositor.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/luna-compositor.log</string>
</dict>
</plist>
